# crontab
0 * * * * root source /opt/venv/bin/activate && /opt/venv/bin/python3 /opt/venv/energyzero.py
